# etc/cron.d/ramdisk
# /etc/cron.d/ramdisk crontab entry
# Default: every 2 hours, RAM-copy-to-DISK

SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* */2 * * *  root [ -x /etc/init.d/ramdisk ] && /usr/sbin/invoke-rc.d ramdisk sync > /dev/null
